# local/overlays/hardware_support_layer/resources/services/boost_clocks_and_fan/boost_clocks_and_fan.service
[Unit]
Description=Boost clocks and fan
After=multi-user.target
Requires=multi-user.target

[Service]
ExecStart=/usr/local/bin/boost_clocks_and_fan.sh


[Install]
WantedBy=multi-user.target