# local/overlays/hardware_support_layer/resources/services/spi_start/spi_start.service
[Unit]
Description=Can start
After=remote-fs.target
After=syslog.target
After=multi-user.target
Requires=multi-user.target

[Service]
ExecStart=/usr/local/bin/spi_start.sh

[Install]
WantedBy=multi-user.target
